# consul-template/consul-template.target
[Unit]
Description=Target file for service-discover-template@ service templates

[Install]
WantedBy=multi-user.target
