# service-discover-template/service-discover-template@-legacy.service
# SPDX-FileCopyrightText: 2023-2024 Zextras <https://www.zextras.com>
#
# SPDX-License-Identifier: AGPL-3.0-only

[Unit]
Description=Zextras Service Discover Template (:%i)

[Service]
User=service-discover
Group=service-discover
EnvironmentFile=-/etc/default/service-discover-template
ExecStart=/usr/bin/consul-template -config /etc/zextras/service-discover-template/%i.hcl
ExecReload=/usr/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
